# .config/supervisord.conf
[supervisord]
user = root
logfile = /var/log/supervisor.log

[program:nginx]
command = nginx

[program:gunicorn]
;배포 시 수정 필요
command = gunicorn -c /srv/amantha/.config/gunicorn.py config.wsgi.dev_hj